# docker-root/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

/usr/sbin/cron -f -L 1
